# deploy/control-plane/systemd/stability-backend-migrate.service
[Unit]
Description=Stability Test Platform Backend DB Migration (oneshot)
After=network.target

[Service]
Type=oneshot
User=<deploy-user>
WorkingDirectory=/home/debian13/stability-test-platform
EnvironmentFile=/home/debian13/stability-test-platform/.env.backend
ExecStart=/bin/bash -lc 'cd /home/debian13/stability-test-platform/backend && ../venv/bin/python -m alembic upgrade head'

[Install]
WantedBy=multi-user.target

