# deploy/control-plane/systemd/stability-backend.service
[Unit]
Description=Stability Test Platform Backend
After=network.target

[Service]
Type=simple
# Replace with your deploy user, e.g. User=stability
User=<deploy-user>
WorkingDirectory=/home/debian13/stability-test-platform
EnvironmentFile=/home/debian13/stability-test-platform/.env.backend
# Keep schema in sync with the running code path before the single backend instance starts.
ExecStartPre=/bin/bash -lc 'cd /home/debian13/stability-test-platform/backend && ../venv/bin/python -m alembic upgrade head'
ExecStart=/home/debian13/stability-test-platform/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=append:/home/debian13/stability-test-platform/logs/backend.log
StandardError=append:/home/debian13/stability-test-platform/logs/backend_error.log

[Install]
WantedBy=multi-user.target
